# ops/systemd/illuvrse-ops.timer
[Unit]
Description=Run ILLUVRSE Ops Orchestrator every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Unit=illuvrse-ops.service

[Install]
WantedBy=timers.target
